# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server

[program:app]
command=node dist/index.js
directory=/app
autostart=true
autorestart=true
